# nginx-fpm/configs/wordpress.cron
*/15 * * * *   curl -s 'http://127.0.0.1:8080/wp-cron.php?doing_wp_cron' > /proc/$(cat /var/run/crond.pid)/fd/1 2>&1
#